# scripts/deploy/ocms.service
# Copyright (c) 2025-2026 Oleg Ivanchenko
# SPDX-License-Identifier: GPL-3.0-or-later
#
# oCMS Systemd Service File
# Copy to: /etc/systemd/system/ocms.service
#
# After copying, run:
#   sudo systemctl daemon-reload
#   sudo systemctl enable ocms
#   sudo systemctl start ocms

[Unit]
Description=oCMS Content Management System
Documentation=https://github.com/ocms-go/ocms
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ocms
Group=ocms
WorkingDirectory=/var/lib/ocms

# Environment configuration
EnvironmentFile=/etc/ocms/ocms.env

# Application binary
ExecStart=/opt/ocms/bin/ocms

# Restart policy
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true

# Allow write access to data directories
ReadWritePaths=/var/lib/ocms

# Resource limits (adjust as needed)
# LimitNOFILE=65535
# MemoryMax=512M
# CPUQuota=100%

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ocms

[Install]
WantedBy=multi-user.target
